# files/rmqttd.in
#!/bin/sh

# PROVIDE: rmqttd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable rmqttd:
#
# rmqttd_enable="YES"
# rmqttd_flags="<set as needed>"
#
# See rmqttd(8) for flags
#

. /etc/rc.subr

name=rmqttd
rcvar=rmqttd_enable

load_rc_config $name

rmqttd_enable=${rmqttd_enable:="NO"}
rmqttd_user=${rmqttd_user:="nobody"}

command=%%PREFIX%%/bin/rmqttd
pidfile=${rmqttd_pidfile:-"/var/run/rmqttd.pid"}

extra_commands=reload
start_precmd=rmqttd_precmd
stop_postcmd="rm -f $pidfile"

rmqttd_precmd()
{
    /usr/bin/install -o ${rmqttd_user} -m 644 /dev/null ${pidfile}
    /usr/bin/install -d -o ${rmqttd_user} -m 755 /var/log/rmqtt
}

run_rc_command "$1"
